# usr/share/init/cgi-keygen.conf
author "Development Infrastructure Team"
description "cgi-keygen"

start on runlevel 1

task

env RANDFILE="/dev/urandom"

script
    $(. /usr/syno/bin/cgi-keygen)
end script
