# usr/share/init/ddnsd.conf
description "DDNS Daemon"

author "Network Infrastructure Team"

start on syno.network.ready
stop on runlevel [06]

console log
expect fork
respawn
respawn limit 5 10
kill timeout 15

exec /usr/syno/sbin/ddnsd

# vim:ft=upstart
