# usr/share/init/DiskHealthCk.conf
description "Schedule disk health check"
author "Business Platform Team"

start on syno.scemd.ready
console log

task

script
	/usr/syno/etc/rc.sysv/DiskHealthCk.sh || true
end script

