# usr/share/init/findhostd.conf
author "System Integration Team"
description	"findhost daemon"

start on syno.network.ready
stop on all-legacy-service-down

oom score -999

expect fork

respawn
respawn limit 5 10

console log

pre-start script
    echo `/bin/date -u`": start findhostd" || true
end script

script
	exec /usr/syno/bin/findhostd
end script

pre-stop script
    echo `/bin/date -u`": stop findhostd" || true
end script
# vim:ft=upstart
