# usr/share/init/hibernationdebug.conf
description "Hibernation debug check"
author "Storage Platform Team"

start on runlevel 1
console log

task

script
	/usr/syno/sbin/syno_hibernation_debug || true
end script

