# usr/share/init/iscsi_pluginengined.conf
description "iSCSI Plugin Daemon"
author "Storage Technology Team"
start on syno.iscsi.ready
stop on stopping iscsitrg
console log
respawn
respawn limit 5 10

script
	exec /usr/syno/bin/iscsi_snapshot_comm_core -D
end script
