# usr/share/init/SmartDbUpdateCk.conf
description "Schedule smart db update"
author "Storage Platform Team"

start on runlevel 2
console log

task

script
	/usr/syno/etc/rc.sysv/SmartDbUpdateCk.sh || true
end script
