# usr/share/init/SynoMemCheck.conf
description "Memory Rule Check"
author "Storage Platform Team"

start on runlevel 2
console log

task

script
	/usr/syno/bin/syno_mem_check || true
end script
