# usr/share/init/synosupportchannelcheck.conf
description "Support channel checker"
author "Web Infrastructure Team"

start on runlevel 2

console log

task

exec /usr/syno/bin/synosupportchannelchecker
