# usr/share/init/taskmanager.conf
description "taskmanager"
author "Web Infrastructure Team"

start on stopped root-file-system

task
console log

script  

	/bin/mkdir -m 1777 /tmp/DSMTaskMgr || true

end script

