# usr/share/init/tty.conf
description	"start getty for console output"

start on stopped rc

respawn
respawn limit 5 10

console none

exec /sbin/getty 115200 console
# vim:ft=upstart
