# usr/share/init/upnpd.conf
description "start upnpd"
author "Network Infrastructure Team"

console log
start on syno.network.ready
stop on runlevel [06]

respawn
respawn limit 5 10
expect fork

pre-start script
        if [ -x "/usr/bin/upnpd" ] ; then
            echo "Starting upnpd..."
        fi 
end script

exec /usr/bin/upnpd

# vim:ft=upstart
