# usr/share/init/WDidle3Dis.conf
description "WD idle3 timer"
author "Business Platform Team"

start on runlevel 1
stop on runlevel [06]


console log

post-stop script
	/usr/syno/etc/rc.sysv/WDidle3Dis.sh stop || true
end script

