# usr/lib/udev/rules.d/50-boot-flash.rules
# Author: Storage Platform Team

SUBSYSTEM=="block|mtd", KERNEL=="synoboot2|mtd0", ACTION=="add", RUN+="/usr/lib/udev/script/syno_flash_ready.sh %k"
